# arptables/files/rc.arptables
#!/bin/bash
#
# init script for arptables
#
# Original by Dag Wieers <dag@wieers.com>.
# Modified/changed to arptables by
# 	Rok Papez <rok.papez@arnes.si>.
#
# chkconfig: - 16 84
# description: Arp filtering tables
#
# config: /etc/sysconfig/arptables

source /etc/init.d/functions

[ -x /usr/sbin/arptables ] || exit 1
[ -x /usr/sbin/arptables-save ] || exit 1
[ -x /usr/sbin/arptables-restore ] || exit 1

[ "$1" != "save" -o -r /etc/sysconfig/arptables ] || exit 1

RETVAL=0
prog="arptables"
desc="Arp filtering"

start() {
	echo -n $"Starting $desc ($prog): "
	/usr/sbin/arptables-restore < /etc/sysconfig/arptables || RETVAL=1

	if [ $RETVAL -eq 0 ]; then
		echo "OK"
		rm -f /var/lock/subsys/$prog
	else
		echo "ERROR"
	fi

	echo
	return $RETVAL
}

stop() {
	echo -n $"Stopping $desc ($prog): "
	/usr/sbin/arptables-restore < /dev/null || RETVAL=1

	if [ $RETVAL -eq 0 ]; then
		echo "OK"
		rm -f %{_localstatedir}/lock/subsys/$prog
	else
		echo "ERROR"
	fi

	echo
	return $RETVAL
}

restart() {
	stop
	start
}

save() {
	echo -n $"Saving $desc ($prog): "
	/usr/sbin/arptables-save > /etc/sysconfig/arptables || RETVAL=1

	if [ $RETVAL -eq 0 ]; then
		echo "OK"
	else
		echo "ERROR"
	fi
	echo
}

case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart|reload)
	restart
	;;
  condrestart)
	[ -e /var/lock/subsys/$prog ] && restart
	RETVAL=$?
	;;
  save)
	save
	;;
  status)
	/usr/sbin/arptables-save
	RETVAL=$?
	;;
  *)
	echo $"Usage $0 {start|stop|restart|condrestart|save|status}"
	RETVAL=1
esac

exit $RETVAL
